# gorouter/templates/default/supervisord.conf.erb
[unix_http_server]
file=/var/run/supervisor.sock   ; (the path to the socket file)
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

;[inet_http_server]         ; inet (TCP) server disabled by default
;port=127.0.0.1:9001        ; (ip_address:port specifier, *:port for all iface)
;username=user              ; (default is no username (open server))
;password=123               ; (default is no password (open server))

[supervisord]
logfile=<%= node['metrichill']['gorouter']['supervisord']['logfile'] %> ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=<%= node['metrichill']['gorouter']['supervisord']['logfile_maxbytes'] %>        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=<%= node['metrichill']['gorouter']['supervisord']['logfile_backups'] %>           ; (num of main logfile rotation backups;default 10)
loglevel=<%= node['metrichill']['gorouter']['supervisord']['loglevel'] %>               ; (log level;default info; others: debug,warn,trace)
pidfile=<%= node['metrichill']['gorouter']['supervisord']['pidfile'] %> ; (supervisord pidfile;default supervisord.pid)
nodaemon=<%= node['metrichill']['gorouter']['supervisord']['nodaemon'] %>               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;environment=KEY="value"     ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:9001 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as http_username if set
;password=123                ; should be same as http_password if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:gnatsd]
command=/usr/bin/gnatsd
priority=<%= node['metrichill']['gorouter']['supervisord']['gnatsd']['priority'] %>
stdout_logfile=<%= node['metrichill']['gorouter']['supervisord']['gnatsd']['access_log'] %>
stdout_logfile_maxbytes=<%= node['metrichill']['gorouter']['supervisord']['gnatsd']['access_log_size'] %>
stdout_logfile_backups=<%= node['metrichill']['gorouter']['supervisord']['gnatsd']['access_log_backups'] %>
stderr_logfile=<%= node['metrichill']['gorouter']['supervisord']['gnatsd']['error_log'] %>
stderr_logfile_maxbytes=<%= node['metrichill']['gorouter']['supervisord']['gnatsd']['error_log_size'] %>
stderr_logfile_backups=<%= node['metrichill']['gorouter']['supervisord']['gnatsd']['error_log_backups'] %>

[program:gorouter]
command=/srv/gorouter/bin/gorouter -c /srv/gorouter/etc/gorouter.yml
priority=<%= node['metrichill']['gorouter']['supervisord']['gorouter']['priority'] %>
stdout_logfile=<%= node['metrichill']['gorouter']['supervisord']['gorouter']['access_log'] %>
stdout_logfile_maxbytes=<%= node['metrichill']['gorouter']['supervisord']['gorouter']['access_log_size'] %>
stdout_logfile_backups=<%= node['metrichill']['gorouter']['supervisord']['gorouter']['access_log_backups'] %>
stderr_logfile=<%= node['metrichill']['gorouter']['supervisord']['gorouter']['error_log'] %>
stderr_logfile_maxbytes=<%= node['metrichill']['gorouter']['supervisord']['gorouter']['error_log_size'] %>
stderr_logfile_backups=<%= node['metrichill']['gorouter']['supervisord']['gorouter']['error_log_backups'] %>
